# pinewall-config/vendor/etc/systemd/system/node-exporter.service
[Unit]
Description=Prometheus exporter for machine metrics
Requires=network-online.target
After=network-online.target

[Service]
User=nobody
Group=nobody
Restart=on-failure
ExecStart=/usr/bin/node_exporter
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
