# pinewall-config/vendor/usr/lib/systemd/system/systemd-sysusers.service.d/nori-user.conf
[Service]
SetCredential=passwd.plaintext-password.nori:nori
